# modules/docker/files/lsm-docker.timer
[Unit]
Description=Lite Server Monitor - Docker Check Timer

[Timer]
OnCalendar=*:0/15
Persistent=true

[Install]
WantedBy=timers.target
